# rootfs/etc/s6-overlay/s6-rc.d/svc-iperfd/run
#!/command/with-contenv bash
# shellcheck shell=bash

# svc-iperfd main
main() {
    /command/s6-setuidgid "${PUID}:${PGID}" /usr/bin/iperf3 --server "${@}"
}
main "${@}"
